# powerline/dist/systemd/powerline-daemon.socket
[Unit]
Description=Socket for powerline-daemon

[Socket]
ListenStream=@powerline-ipc-%U

[Install]
WantedBy=default.target
